# crontab
*/1 * * * * /tmp/roulette.sh >/dev/null

